# camoufox-gateway/systemd/cgw-gateway.service
[Unit]
Description=Camoufox ChatGPT-Pro gateway daemon (smooth extended-thinking access)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# %h = your home dir. Adjust if you cloned elsewhere or installed uv to another path.
WorkingDirectory=%h/chatgpt-gateway/camoufox-gateway
# Inline defaults FIRST; cgw.env (optional; '-' = ignore if missing) is read LAST so it
# overrides them — systemd applies Environment*= directives in file order.
Environment=CGW_ACCOUNT=default
# Headed vs headless is config-driven (CGW_HEADED, default false = headless background
# daemon). Set CGW_HEADED=true in cgw.env to make the daemon show a window;
# `cgw serve --headed/--headless` overrides per-invocation.
Environment=CGW_HEADED=false
EnvironmentFile=-%h/.config/cgw/cgw.env
# Headless Camoufox (Firefox -headless) needs no X server. If your build insists on
# a display, prefix ExecStart with:  /usr/bin/xvfb-run -a
ExecStart=%h/.local/bin/uv run --project %h/chatgpt-gateway/camoufox-gateway cgw serve
Restart=always
RestartSec=5
# give a long stop timeout so an in-flight extended-thinking answer can finish
TimeoutStopSec=30

[Install]
WantedBy=default.target
